# files/system/usr/lib/systemd/system/birdie-system-setup.service
[Unit]
Description=Configure birdie system
After=rpm-ostreed.service tuned.service
Before=systemd-user-sessions.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/birdie-system-setup

[Install]
WantedBy=multi-user.target